# systemd/ice-gateway.service
[Unit]
Description=Ice Gateway Monitor
Documentation=https://github.com/your-repo/scotsman-monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=APP_USER
WorkingDirectory=/opt/ice_gateway
EnvironmentFile=/etc/ice-gateway/ice-gateway.env
Environment="PATH=/home/APP_USER/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/APP_USER/.local/bin/uv run ice-gateway
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ice-gateway

[Install]
WantedBy=multi-user.target
